# meta-ibm/recipes-phosphor/ips-asset-tag/ips-asset-tag/ips-asset-tag.path
[Unit]
Description=Monitor the AssetTag file for IPS

[Path]
PathModified=/var/lib/phosphor-inventory-manager/xyz/openbmc_project/inventory/system/xyz.openbmc_project.Inventory.Decorator.AssetTag
Unit=ips-copy-asset-tag.service

[Install]
WantedBy=multi-user.target
