# meta-ibm/recipes-phosphor/ips-asset-tag/ips-asset-tag/ips-update-asset-tag.service
[Unit]
Description=Copy the AssetTag file for IPS

[Service]
ExecStart=/usr/bin/ips-asset-tag.sh update
