# meta-ibm/meta-witherspoon/recipes-phosphor/gpio/id-button/99-gpio-keys.rules
ACTION=="add", SUBSYSTEM=="input", ENV{ID_PATH}=="platform-gpio-keys", RUN+="/bin/systemctl start phosphor-gpio-monitor@id_button.service"
